# ami/git/git-credential-refresh.service
# User systemd unit: keep git GitHub App credentials refreshed in the background.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp /path/to/ami/git/git-credential-refresh.service ~/.config/systemd/user/
#   # Edit ExecStart and Environment=GITHUB_KEY_SECRET_CONFIG if needed.
#   systemctl --user daemon-reload
#   systemctl --user enable --now git-credential-refresh.service
#
# Logs:
#   journalctl --user -u git-credential-refresh.service -f

[Unit]
Description=Refresh git GitHub App credentials from AWS Secrets Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=GITHUB_KEY_SECRET_CONFIG=/home/ubuntu/.config/git-auth/bots.json
# Environment=GIT_AUTH_REFRESH_INTERVAL_SECONDS=1800
# Environment=AWS_REGION=us-east-1
# Environment=AWS_PROFILE=default

ExecStart=/path/to/ami/git/git-credential-refresh-daemon.sh
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
