# ami/git/git-credential-refresh.service
# Systemd user unit: run the git credential refresh daemon on startup and keep it running.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp git-credential-refresh.service ~/.config/systemd/user/
#   # Edit the unit to set your secret name, App ID, and Installation ID (or use a drop-in).
#   systemctl --user daemon-reload
#   systemctl --user enable --now git-credential-refresh.service
#
# Or install for all users (system-wide): copy to /etc/systemd/system/ and set User=.
# Then: systemctl enable --now git-credential-refresh.service
#
# Required: set the three env vars (e.g. in this file or via environment.d drop-in).

[Unit]
Description=Git credential refresh daemon (GitHub App token)
After=network-online.target

[Service]
Type=simple
# Set these for your GitHub App (or use EnvironmentFile=).
Environment=GITHUB_KEY_SECRET_NAME=
Environment=GITHUB_APP_ID=
Environment=GITHUB_INSTALLATION_ID=
# Optional: custom token file path and refresh interval (default 30 min).
# Environment=GIT_AUTH_TOKEN_FILE=%h/.config/git-auth/github-token
# Environment=GIT_AUTH_REFRESH_INTERVAL_SECONDS=1800

ExecStart=/bin/bash -c 'exec /path/to/ami/git/git-credential-refresh-daemon.sh'
# Replace /path/to/ami with the actual path to the repo (e.g. /home/ubuntu/desk), so the script path is ami/git/git-credential-refresh-daemon.sh

Restart=always
RestartSec=60

[Install]
WantedBy=default.target
