# ami/git/git-credential-refresh.service
# User systemd unit: keep git GitHub App credentials refreshed in the background.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp /path/to/ami/git/git-credential-refresh.service ~/.config/systemd/user/
#   # Edit the unit and set GITHUB_KEY_SECRET_NAME (and optional env) and ExecStart path.
#   systemctl --user daemon-reload
#   systemctl --user enable --now git-credential-refresh.service
#
# Logs:
#   journalctl --user -u git-credential-refresh.service -f

[Unit]
Description=Refresh git GitHub App credentials from AWS Secrets Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Set your secret name and optional AWS/git-auth env. Example:
Environment=GITHUB_KEY_SECRET_NAME=my-github-app
# Environment=GIT_AUTH_REFRESH_INTERVAL_SECONDS=1800
# Environment=AWS_REGION=us-east-1
# Environment=AWS_PROFILE=default

# Path to the daemon script (edit to match your ami/git location)
ExecStart=/path/to/ami/git/git-credential-refresh-daemon.sh
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
